# tools/guidashboard-upgrade.service
[Unit]
Description=GUI Dashboard Upgrade (runs deploy + install)
After=network.target

[Service]
Type=oneshot
ExecStartPre=/bin/sh -c 'echo "=== Upgrade started $(date) ===" > /tmp/upgrade.log && chmod 644 /tmp/upgrade.log'
ExecStart=/usr/local/bin/run-upgrade
StandardOutput=append:/tmp/upgrade.log
StandardError=append:/tmp/upgrade.log
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target
